# docker/build/supervisor/conf.d/getmail-poll.conf
# Processes configs via a single getmail command sequentially in a batch (implicit):
[program:getmail-poll]
#stdout_logfile=/var/log/supervisor/%(program_name)s.log
#stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart = true
killasgroup = true
command = bash -c 'getmail --getmaildir=/var/lib/getmail --rcfile=/run/getmail/configs/poll/jane.doe@retriever.test && sleep %(ENV_GETMAIL_POLL)s'
