# docker/build/supervisor/conf.d/getmail-idle.conf
# Run one getmail command per IDLE config:
[program:getmail-idle]
numprocs = 1
numprocs_start = 1
process_name = %(program_name)s_%(process_num)s
#stdout_logfile=/var/log/supervisor/%(program_name)s_%(process_num)s.log
#stderr_logfile=/var/log/supervisor/%(program_name)s_%(process_num)s.log
autorestart = true
killasgroup = true
command = getmail
  --getmaildir=/var/lib/getmail
  --rcfile=/run/getmail/configs/idle/john.doe@remote.test.cf
  --idle=INBOX
